# meta-unit-graphics/dynamic-layers/flutter-layer/recipes-core/systemd/systemd-flutter-kiosk/kiosk.path
[Unit]
Description=Wait for weston

[Path]
PathExists=%h/.local/share/weston/weston-ready-signal

[Install]
WantedBy=graphical.target
