# meta-unit-core/dynamic-layers/virtualization-layer/recipes-core/systemd/systemd-distribution-path/distribution.path
[Unit]
Description=Wait for TLS cert and key

[Path]
PathExists=%h/.local/share/distribution/certs-ready-signal

[Install]
WantedBy=default.target
